# docker/supervisord.conf
; Lumimax 单镜像 supervisord 配置
; 负责守护 nginx + gateway + base-service + biz-service + iot-service
;
; 日志策略：所有进程的 stdout/stderr 直接转到容器 stdout/stderr，
; 由 docker logs 接管。NestJS 内部 pino 仍会写文件到 /app/api/logs/*.log。

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/app/run/supervisord.pid
loglevel=info
silent=false

[unix_http_server]
file=/app/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///app/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ===================== nginx =====================
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10
startsecs=3
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ===================== gateway =====================
[program:gateway]
command=node /app/api/dist/apps/gateway/src/main.js
directory=/app/api/apps/gateway
user=app
autostart=true
autorestart=true
priority=18
startsecs=5
stopwaitsecs=15
environment=
    NODE_ENV="production",
    NODE_PATH="/app/api/node_modules/.pnpm/node_modules:/app/api/apps/gateway/node_modules:/app/api/node_modules",
    SERVICE_NAME="gateway",
    HOST="127.0.0.1",
    HTTP_PORT="4000",
    BASE_SERVICE_GRPC_URL="127.0.0.1:4120",
    BIZ_SERVICE_GRPC_URL="127.0.0.1:4130",
    IOT_SERVICE_GRPC_ENDPOINT="127.0.0.1:4140"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ===================== base-service =====================
; 优先级低于 nginx、先于 biz：先跑迁移/占 advisory lock，减轻同时连库导致的瞬时断连
[program:base-service]
command=node /app/api/dist/apps/base-service/src/main.js
directory=/app/api/apps/base-service
user=app
autostart=true
autorestart=true
priority=14
startsecs=5
stopwaitsecs=15
environment=
    NODE_ENV="production",
    NODE_PATH="/app/api/node_modules/.pnpm/node_modules:/app/api/apps/base-service/node_modules:/app/api/node_modules",
    SERVICE_NAME="base-service",
    HOST="127.0.0.1",
    HTTP_PORT="4020",
    GRPC_PORT="4120"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ===================== biz-service =====================
[program:biz-service]
command=node /app/api/dist/apps/biz-service/src/main.js
directory=/app/api/apps/biz-service
user=app
autostart=true
autorestart=true
priority=24
startsecs=5
stopwaitsecs=15
environment=
    NODE_ENV="production",
    NODE_PATH="/app/api/node_modules/.pnpm/node_modules:/app/api/apps/biz-service/node_modules:/app/api/node_modules",
    SERVICE_NAME="biz-service",
    HOST="127.0.0.1",
    HTTP_PORT="4030",
    GRPC_PORT="4130"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ===================== iot-service =====================
[program:iot-service]
command=node /app/api/dist/apps/iot-service/src/main.js
directory=/app/api/apps/iot-service
user=app
autostart=true
autorestart=true
priority=22
startsecs=5
stopwaitsecs=15
environment=
    NODE_ENV="production",
    NODE_PATH="/app/api/node_modules/.pnpm/node_modules:/app/api/apps/iot-service/node_modules:/app/api/node_modules",
    SERVICE_NAME="iot-service",
    HOST="127.0.0.1",
    HTTP_PORT="4040",
    GRPC_PORT="4140"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ===================== 分组 =====================
[group:lumimax]
programs=gateway,base-service,biz-service,iot-service,nginx
priority=999
